# admin/crontab
* * * * * node /usr/src/app.js >> /var/log/cron.log 2>&1
